# server/botland.conf
description     "Botland: AI MMO Server"

start on startup

chdir /root/botland
respawn
respawn limit 10 5
exec runghc Api >> /root/botland/log.txt 2>&1