# server/botland.conf
description     "Botland: AI MMO Server"

start on startup

chdir /root/botland
respawn
respawn limit 5 5 
exec ./Api 

# >> /root/botland/log.txt 2>&1